# files/misc/crontab
*/5 * * * * /home/rui/dotfiles/files/scripts/battery-notify
