# files/misc/crontab
*/5 * * * * /home/rui/dotfiles/scripts/battery-notify
